# ansible/roles/tasks-systemd/cockpituous-janitor.timer
[Unit]
Description=Regularly clean up cockpituous cruft

[Timer]
OnBootSec=1h
OnUnitActiveSec=24h
AccuracySec=15min
