# verify/cockpit-verify.service
[Unit]
Description=Cockpit Verify

[Service]
User=cockpit
Group=cockpit
Environment=TEST_OS=fedora-23
RestartSec=60
WorkingDirectory=/opt/cockpit/test
ExecStart=/opt/cockpit-verify
SendSIGKILL=no
SendSIGHUP=yes

[Install]
WantedBy=multi-user.target
